# supervisor_forex_bot.conf
; ============================================================
; ABSOLUTE PATH: /etc/supervisor/conf.d/forex_bot.conf
; Supervisor process management for all forex bot services
; Install: sudo cp forex_bot.conf /etc/supervisor/conf.d/
;          sudo supervisorctl reread
;          sudo supervisorctl update
; ============================================================

; ── Django / Daphne ASGI server ──────────────────────────────
[program:forex_bot_web]
command=/opt/forex_bot_venv/bin/daphne
    -b 127.0.0.1
    -p 8001
    --websocket_timeout 86400
    config.asgi:application
directory=/opt/forex_bot
user=forex
environment=DJANGO_SETTINGS_MODULE="config.settings"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/forex_bot/logs/daphne.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopasgroup=true
killasgroup=true
priority=100

; ── Celery Worker — Trading queue (bot loops) ────────────────
[program:celery_trading]
command=/opt/forex_bot_venv/bin/celery
    -A config.celery worker
    --queues=trading
    --concurrency=4
    --prefetch-multiplier=1
    --max-tasks-per-child=50
    --loglevel=info
    --logfile=/opt/forex_bot/logs/celery_trading.log
    --hostname=trading@%%h
directory=/opt/forex_bot
user=forex
environment=DJANGO_SETTINGS_MODULE="config.settings"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/forex_bot/logs/celery_trading.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=120
stopasgroup=true
killasgroup=true
priority=200

; ── Celery Worker — Orders queue (fast order execution) ──────
[program:celery_orders]
command=/opt/forex_bot_venv/bin/celery
    -A config.celery worker
    --queues=orders,commands
    --concurrency=8
    --prefetch-multiplier=1
    --max-tasks-per-child=200
    --loglevel=info
    --hostname=orders@%%h
directory=/opt/forex_bot
user=forex
environment=DJANGO_SETTINGS_MODULE="config.settings"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/forex_bot/logs/celery_orders.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
priority=200

; ── Celery Worker — Data + Default queues ────────────────────
[program:celery_data]
command=/opt/forex_bot_venv/bin/celery
    -A config.celery worker
    --queues=data,default
    --concurrency=4
    --prefetch-multiplier=2
    --max-tasks-per-child=100
    --loglevel=info
    --hostname=data@%%h
directory=/opt/forex_bot
user=forex
environment=DJANGO_SETTINGS_MODULE="config.settings"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/forex_bot/logs/celery_data.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopasgroup=true
killasgroup=true
priority=200

; ── Celery Worker — Backtesting queue (CPU-intensive) ────────
[program:celery_backtesting]
command=/opt/forex_bot_venv/bin/celery
    -A config.celery worker
    --queues=backtesting
    --concurrency=2
    --prefetch-multiplier=1
    --max-tasks-per-child=10
    --loglevel=info
    --hostname=backtesting@%%h
directory=/opt/forex_bot
user=forex
environment=DJANGO_SETTINGS_MODULE="config.settings"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/forex_bot/logs/celery_backtesting.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=300
stopasgroup=true
killasgroup=true
priority=300

; ── Celery Beat (periodic task scheduler) ────────────────────
[program:celery_beat]
command=/opt/forex_bot_venv/bin/celery
    -A config.celery beat
    --scheduler=django_celery_beat.schedulers:DatabaseScheduler
    --loglevel=info
    --logfile=/opt/forex_bot/logs/celery_beat.log
    --pidfile=/tmp/celerybeat.pid
directory=/opt/forex_bot
user=forex
environment=DJANGO_SETTINGS_MODULE="config.settings"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/forex_bot/logs/celery_beat.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
priority=400

; ── Group: all celery workers ────────────────────────────────
[group:celery]
programs=celery_trading,celery_orders,celery_data,celery_backtesting,celery_beat
priority=200
